# com.atlas.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.atlas.menubar</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_PATH__</string>
        <string>__INSTALL_DIR__/atlas/menubar_agent.py</string>
        <string>--agent-port</string>
        <string>8767</string>
        <string>--fleet-server</string>
        <string>https://localhost:__FLEET_PORT__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__ACTUAL_HOME__/.atlas-menubar.log</string>

    <key>StandardErrorPath</key>
    <string>__ACTUAL_HOME__/.atlas-menubar.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
