# com.atlas.fleetserver.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.atlas.fleetserver</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_PATH__</string>
        <string>-m</string>
        <string>atlas.fleet_server</string>
        <string>--port</string>
        <string>__FLEET_PORT__</string>
        <string>--host</string>
        <string>0.0.0.0</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/var/log/atlas-fleetserver.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/atlas-fleetserver.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <key>UserName</key>
    <string>__ACTUAL_USER__</string>

    <key>GroupName</key>
    <string>staff</string>
</dict>
</plist>
